# crontab
* * * * * prefect /app/scripts/bitbucket-fetcher.sh
* * * * * prefect /app/scripts/gitlab-fetcher.sh