# crontab
0 3 * * * prefect /app/scripts/bitbucket-fetcher.sh
0 4 * * * prefect /app/scripts/gitlab-fetcher.sh