# crontab
# run every Friday at 11:59 PM as user “prefect”
59 23 * * 5 prefect /app/scripts/bitbucket-fetcher.sh
59 23 * * 5 prefect /app/scripts/gitlab-fetcher.sh
