# crontab
0 3 * * * /app/scripts/bitbucket-fetcher.sh
0 4 * * * /app/scripts/gitlab-fetcher.sh